# utils/edge-utils/files/etc/init.d/btuart
#!/bin/sh /etc/rc.common
# btuart init script

START=99

start() {
	hciconfig hci0 up
	btuart
}