# utils/edge-utils/files/etc/uci-defaults/99_openssh-permit-root-login
if [ -f /etc/ssh/sshd_config ]; then
	mv /etc/ssh/sshd_config /etc/ssh/sshd_config.orig
	cat /etc/ssh/sshd_config.orig | sed 's/#PermitRootLogin prohibit-password/PermitRootLogin yes/g' > /etc/ssh/sshd_config
fi

exit 0